# systemd/easy-portainer-cleanup.timer
[Unit]
Description=Garbage collect docker registry periodically

[Timer]
OnStartupSec=5min
OnUnitInactiveSec=5min

[Install]
WantedBy=timers.target
